# a10_octavia/bin/a10-controller-worker.service
[Unit]
Description = Devstack devstack@o-a10cw.service

[Service]
ExecReload = /bin/kill -HUP $MAINPID
TimeoutStopSec = 300
KillMode = process
ExecStart = /usr/local/bin/a10-octavia-worker  --config-file /etc/octavia/octavia.conf
User = stack

[Install]
WantedBy = multi-user.target
